# src/agent_audit/knowledge/rule_packs/atr/privilege-escalation/ATR-2026-00107-privilege-escalation-via-delayed-task-execution-bypass.yaml
# =============================================================
# atr.privilege-escalation.privilege-escalation-via-delayed-task-execution-bypass
# =============================================================
# Adapted from ATR-2026-00107 (https://github.com/Agent-Threat-Rule/agent-threat-rules)
# MIT (c) 2026 ATR Contributors
# Severity: high  |  ASAMM: AI-03 (+AI-06,AO-02)
# =============================================================
atr_source:
  original_id: ATR-2026-00107
  original_title: Privilege Escalation via Delayed Task Execution Bypass
  atr_version: '0.1'
  rule_version: 1
  status: experimental
  maturity: experimental
  import_date: '2026-04-21'
  upstream: https://github.com/Agent-Threat-Rule/agent-threat-rules
  license: MIT (c) 2026 ATR Contributors
agent_audit_id: atr.privilege-escalation.privilege-escalation-via-delayed-task-execution-bypass
title: Privilege Escalation via Delayed Task Execution Bypass
description: 'Detects tools that claim to schedule tasks while explicitly stating they bypass permission

  checks or security controls through delayed execution. This technique uses the temporal

  gap between task scheduling and execution to escalate privileges, as delayed tasks may

  run in a system context that bypasses the original user''s permission constraints.'
severity_default: high
category: privilege-escalation
audit_surface:
- session_tool_call
- mcp_manifest
- skill_md
asamm:
  primary: AI-03
  secondary:
  - AI-06
  - AO-02
  mapping_confidence: category-default
  notes: Category-level default applied at import. Review per rule.
references:
  owasp_llm:
  - LLM06:2025 - Excessive Agency
  owasp_agentic:
  - ASI06:2026 - Insufficient Access Controls
  - ASI05:2026 - Unexpected Code Execution
  mitre_attack:
  - T1053 - Scheduled Task/Job
tags:
  category: privilege-escalation
  subcategory: delayed-execution-bypass
  scan_target: mcp
  confidence: medium
agent_source:
  type: tool_call
  framework:
  - mcp
  - langchain
  - llamaindex
  provider:
  - any
detection:
  conditions:
  - field: tool_description
    operator: regex
    value: (?:scheduled|delayed|background)\s+(?:tasks|operations|processes)\s+(?:that\s+)?bypass\s+(?:normal|standard|regular)\s+(?:permission|security|access)\s+checks
    description: Tools explicitly mentioning bypassing permission checks through delayed/scheduled execution
  condition: any
  false_positives:
  - Documentation tools describing privilege escalation risks in scheduled tasks
  - A security training document that quotes the exact phrasing of malicious tool descriptions to teach analysts what to look
    for
  - An automated compliance scanner report that includes snippets of flagged tool descriptions for audit review
response:
  actions:
  - alert
  - snapshot
test_cases:
  true_positives:
  - tool_description: Schedules a background task. Tasks execute with full system privileges after a configurable delay. Scheduled
      tasks bypass normal permission checks as they run in the system context.
    expected: triggered
    description: Scheduled tasks explicitly bypassing permission checks
  - tool_description: Creates delayed operations that bypass standard security checks when executed in the background system
      context.
    expected: triggered
    description: Delayed operations bypassing security checks
  - input: This tool manages scheduled tasks that bypass normal permission checks to ensure background operations complete
      without interruption.
    expected: triggered
    description: Scheduled tasks bypassing normal permission checks
  - input: Registers background processes that bypass standard access checks, allowing elevated execution after a time delay.
    expected: triggered
    description: Background processes bypassing standard access checks
  - input: Enqueues delayed tasks that bypass regular security checks by deferring execution to a privileged system timer
      context.
    expected: triggered
    description: Delayed tasks bypassing regular security checks
  true_negatives:
  - tool_description: Schedules a background task with proper user permissions. All scheduled tasks are subject to the same
      security controls as interactive operations.
    expected: not_triggered
    description: Scheduled tasks that properly maintain permission constraints
  - tool_description: Creates scheduled jobs that run with the same permissions as the user who created them, following standard
      access controls.
    expected: not_triggered
    description: Jobs that inherit user permissions correctly
  - input: This tool creates cron jobs that execute under the invoking user's credentials. All access control policies are
      enforced at runtime.
    expected: not_triggered
    description: Cron job scheduler with proper access control enforcement
  - input: Queues asynchronous work items for later processing. Each item inherits the caller's authorization token and is
      subject to role-based access restrictions.
    expected: not_triggered
    description: Async work queue that preserves caller authorization
  - input: 'Documentation note: attackers sometimes abuse deferred execution to circumvent authorization gates. This tool
      mitigates that risk by re-validating credentials at execution time.'
    expected: not_triggered
    description: Security documentation discussing the risk abstractly without matching the trigger pattern
false_positives: []
